# etc/xdg/autostart/calamares-arcolinux.desktop
[Desktop Entry]
Type=Application
Name=Install ArcoLinuxD
Comment=Installer for ArcoLinuxD
Categories=Utility;
Exec=pkexec /usr/bin/calamares
Icon=arco-calamares
